# fs-overlay/etc/network/if-up.d/info
#!/bin/sh

[ "$IFACE" != "lo" ] || exit 0

/sbin/ifconfig > /dev/tty1
